# config/backup_nettoyage.service
[Unit]
Description=Nettoyage et sauvegarde
After=network-online.target
Wants=network-online.target
ConditionFileIsExecutable=/home/arnaud/Documents/Linux/backup_nettoyage.sh

[Service]
Type=oneshot
ExecStart=/home/arnaud/Documents/Linux/backup_nettoyage.sh
